# presence-monitor/rootfs/etc/services.d/dbus/run
#!/usr/bin/env bashio
# ==============================================================================
# Limych's Hass.io Add-ons: Bluetooth Presence Monitor
# Starts the D-Bus daemon
# ==============================================================================

exec dbus-daemon --system --nofork \
  || bashio::exit.nok 'Failed to start D-Bus daemon'
