# presence-monitor/rootfs/etc/services.d/presence/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Limych's Hass.io Add-ons: Bluetooth Presence Monitor
# Run Bluetooth Presence Monitor
# ==============================================================================
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services
